# executors-roles/deploy/templates/cider-ci_executor-download.service
[Unit]
Description=Cider-CI Executor-Download

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/cider-ci/executor/bin/download
